# quality_uploader/quantlab-quality-uploader.service
[Unit]
Description=QuantLab Quality JSON S3 Uploader Service
Documentation=https://github.com/quantlab/collector
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/quantlab/quality_uploader
EnvironmentFile=/opt/quantlab/.env
ExecStart=/usr/bin/python3 /opt/quantlab/quality_uploader/quality_uploader.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quantlab-quality-uploader

# Resource limits
MemoryMax=512M
CPUQuota=30%

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/quantlab/quality /opt/quantlab/quality_uploader

[Install]
WantedBy=multi-user.target
